# over/etc/s6-init/run/s6/uncaught-logs/run
#!/usr/bin/execlineb -P
fdmove 2 0 fdmove -c 1 2
redirfd -r 0 in
s6-envuidgid -D 65534:65534 nobody
s6-applyuidgid -U
s6-log -bp
t /var/log/uncaught
T 1
